# bin/supervisor_gunicorn.ini
[program:majiong]
directory=/home/App/mahjong_calculator/
command=gunicorn run:app
autostart=true
autorestart=true
startretries=3
redirect_stderr=true
stdout_logfile=/data/logs/supervisor/model_gunicorn.log
stopasgroup=true
killasgroup=true

